# dot_local/share/udev/rules.d/10-coretemp-stable-path.rules
ACTION=="add", SUBSYSTEM=="hwmon", ATTRS{name}=="coretemp", RUN+="/bin/dash -c 'ln -s /sys$devpath/temp1_input /dev/pkg_temp'"